# conf/chronograph.cron
* * * * * @APP_OWNER@ @DEPLOY_DIR@/cron/chronograph.sh