# conf/chronograph.cron
* * * * * @APP_OWNER@ @DEPLOY_DIR@/manage.py cron